# company-brain/submit contents/company-brain/systemd/company-brain-ccwatch.service
[Unit]
Description=Company Brain Claude Code watcher
After=network.target company-brain-server.service
Wants=company-brain-server.service

[Service]
Type=simple
User=chr
WorkingDirectory=/srv/company-brain
EnvironmentFile=/srv/company-brain/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/srv/company-brain/.venv/bin/python -m collectors.claude_code_watcher --interval 60
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
